# debian/postinst
#!/bin/sh
set -e

# Source debconf library
. /usr/share/debconf/confmodule

case "$1" in
    configure)
        # Create system user if it doesn't exist
        if ! getent passwd fibratus > /dev/null; then
            adduser --system --home /var/lib/fibratus --group --quiet fibratus
        fi
        
        # Set permissions
        chown -R fibratus:fibratus /etc/fibratus
        chown -R fibratus:fibratus /var/lib/fibratus
        chown -R fibratus:fibratus /usr/share/fibratus
        
        # Run installation script
        if [ -x /usr/share/fibratus/scripts/install.sh ]; then
            echo "Running installation script..."
            /usr/share/fibratus/scripts/install.sh
        fi
        
        # Enable and start service
        systemctl daemon-reload
        systemctl enable fibratus-portal.service
        systemctl start fibratus-portal.service
        
        # Set up NGINX
        if [ -f /etc/nginx/sites-available/fibratus-portal ]; then
            if [ ! -f /etc/nginx/sites-enabled/fibratus-portal ]; then
                ln -sf /etc/nginx/sites-available/fibratus-portal /etc/nginx/sites-enabled/
                systemctl reload nginx
            fi
        fi
        
        echo "Fibratus Management Portal installation completed."
        echo "Access the portal at https://your-server-ip"
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0